# files/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
user = root
autostart = true

[program:php-fpm]
command = /usr/bin/php-fpm -FR
user = root
autostart = true
